# templates/systemd.service.erb
# This file is managed by Puppet. DO NOT EDIT.
# If you modify this, please also make sure to edit init.sh

[Unit]
Description=Open source monitoring and visualization UI for the entire TICK stack.
Documentation="https://www.influxdata.com/time-series-platform/chronograf/"
After=network-online.target

[Service]
User=<%= scope.lookupvar("chronograf::user") %>
Group=<%= scope.lookupvar("chronograf::group") %>
Environment="HOST=<%= scope.lookupvar("chronograf::host") %>"
Environment="PORT=<%= scope.lookupvar("chronograf::port") %>"
Environment="BOLT_PATH=<%= scope.lookupvar("chronograf::bolt_path") %>"
Environment="CANNED_PATH=<%= scope.lookupvar("chronograf::canned_path") %>"
Environment="PROTOBOARDS_PATH=<%= scope.lookupvar("chronograf::protoboards_path") %>"
# Environment="RESOURCES_PATH=<%= scope.lookupvar("chronograf::resources_path") %>"
# Environment="BASE_PATH=<%= scope.lookupvar("chronograf::basepath") %>"
# Environment="$STATUS_FEED_URL=<%= scope.lookupvar("chronograf::status_feed_url") %>"
EnvironmentFile=-<%= scope.lookupvar("chronograf::service_defaults") %>
ExecStart=/usr/bin/chronograf $CHRONOGRAF_OPTS
KillMode=control-group
Restart=on-failure

[Install]
WantedBy=multi-user.target
