# templates/systemd.service.erb
# This file is managed by Puppet. DO NOT EDIT.
# If you modify this, please also make sure to edit init.sh

[Unit]
Description=Open source monitoring and visualization UI for the entire TICK stack.
Documentation=https://docs.influxdata.com/chronograf/
After=network-online.target

[Service]
User=<%= @user %>
Group=<%= @group %>
Environment="HOST=0.0.0.0"
Environment="PORT=8888"
Environment="BOLT_PATH=<%= @bolt_path %>"
Environment="CANNED_PATH=<%= @canned_path %>"
Environment="PROTOBOARDS_PATH=<%= @protoboards_path %>"
Environment="RESOURCES_PATH=<%= @resources_path %>"
# Environment="BASE_PATH=<%= @basepath %>"
# Environment="STATUS_FEED_URL=<%= @status_feed_url %>"
EnvironmentFile=-<%= @service_defaults %>
ExecStart=/usr/bin/chronograf $CHRONOGRAF_OPTS
KillMode=control-group
Restart=on-failure

[Install]
WantedBy=multi-user.target
